# supervisord.conf
[supervisord]
nodaemon=true
user=appuser
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:backend]
command=node /app/backend/src/index.js
directory=/app/backend
user=appuser
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/backend.err.log
stdout_logfile=/var/log/supervisor/backend.out.log
environment=NODE_ENV="production",DATABASE_URL="file:/app/backend/prisma/dev.db",PORT="80"

[program:frontend]
command=node /app/frontend/server.js
directory=/app/frontend
user=appuser
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/frontend.err.log
stdout_logfile=/var/log/supervisor/frontend.out.log
environment=NODE_ENV="production",NEXT_PUBLIC_API_URL="http://localhost:80",PORT="3000",HOSTNAME="0.0.0.0"
depends_on=backend

[eventlistener:processes]
command=bash -c "echo 'READY' && while read line; do echo $line; done"
events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL 